# spotify-kiosk.service
[Unit]
Description=Spotify Kids Player Kiosk Mode
After=graphical.target spotify-player.service
Wants=graphical.target
Requires=spotify-player.service

[Service]
Type=simple
User=spotify-kids
Group=spotify-kids
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/spotify-kids/.Xauthority"
Environment="HOME=/home/spotify-kids"

# Wait for X11 to be ready
ExecStartPre=/bin/bash -c 'until [ -S /tmp/.X11-unix/X0 ]; do sleep 2; done'
ExecStartPre=/bin/sleep 5

# Start kiosk
ExecStart=/opt/spotify-kids/kiosk_launcher.sh

# Restart if crashes
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target